# ports/compat-32/avahi-32/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
drwxr-xr-x	root/root	usr/lib32/
drwxr-xr-x	root/root	usr/lib32/avahi/
-rw-r--r--	root/root	usr/lib32/libavahi-client.a
-rwxr-xr-x	root/root	usr/lib32/libavahi-client.la
lrwxrwxrwx	root/root	usr/lib32/libavahi-client.so -> libavahi-client.so.3.2.9
lrwxrwxrwx	root/root	usr/lib32/libavahi-client.so.3 -> libavahi-client.so.3.2.9
-rwxr-xr-x	root/root	usr/lib32/libavahi-client.so.3.2.9
-rw-r--r--	root/root	usr/lib32/libavahi-common.a
-rwxr-xr-x	root/root	usr/lib32/libavahi-common.la
lrwxrwxrwx	root/root	usr/lib32/libavahi-common.so -> libavahi-common.so.3.5.4
lrwxrwxrwx	root/root	usr/lib32/libavahi-common.so.3 -> libavahi-common.so.3.5.4
-rwxr-xr-x	root/root	usr/lib32/libavahi-common.so.3.5.4
-rw-r--r--	root/root	usr/lib32/libavahi-core.a
-rwxr-xr-x	root/root	usr/lib32/libavahi-core.la
lrwxrwxrwx	root/root	usr/lib32/libavahi-core.so -> libavahi-core.so.7.1.0
lrwxrwxrwx	root/root	usr/lib32/libavahi-core.so.7 -> libavahi-core.so.7.1.0
-rwxr-xr-x	root/root	usr/lib32/libavahi-core.so.7.1.0
-rw-r--r--	root/root	usr/lib32/libavahi-glib.a
-rwxr-xr-x	root/root	usr/lib32/libavahi-glib.la
lrwxrwxrwx	root/root	usr/lib32/libavahi-glib.so -> libavahi-glib.so.1.0.2
lrwxrwxrwx	root/root	usr/lib32/libavahi-glib.so.1 -> libavahi-glib.so.1.0.2
-rwxr-xr-x	root/root	usr/lib32/libavahi-glib.so.1.0.2
-rw-r--r--	root/root	usr/lib32/libavahi-gobject.a
-rwxr-xr-x	root/root	usr/lib32/libavahi-gobject.la
lrwxrwxrwx	root/root	usr/lib32/libavahi-gobject.so -> libavahi-gobject.so.0.0.5
lrwxrwxrwx	root/root	usr/lib32/libavahi-gobject.so.0 -> libavahi-gobject.so.0.0.5
-rwxr-xr-x	root/root	usr/lib32/libavahi-gobject.so.0.0.5
-rw-r--r--	root/root	usr/lib32/libavahi-libevent.a
-rwxr-xr-x	root/root	usr/lib32/libavahi-libevent.la
lrwxrwxrwx	root/root	usr/lib32/libavahi-libevent.so -> libavahi-libevent.so.1.0.0
lrwxrwxrwx	root/root	usr/lib32/libavahi-libevent.so.1 -> libavahi-libevent.so.1.0.0
-rwxr-xr-x	root/root	usr/lib32/libavahi-libevent.so.1.0.0
-rw-r--r--	root/root	usr/lib32/libdns_sd.a
-rwxr-xr-x	root/root	usr/lib32/libdns_sd.la
lrwxrwxrwx	root/root	usr/lib32/libdns_sd.so -> libdns_sd.so.1.0.0
lrwxrwxrwx	root/root	usr/lib32/libdns_sd.so.1 -> libdns_sd.so.1.0.0
-rwxr-xr-x	root/root	usr/lib32/libdns_sd.so.1.0.0
drwxr-xr-x	root/root	usr/lib32/pkgconfig/
-rw-r--r--	root/root	usr/lib32/pkgconfig/avahi-client.pc
-rw-r--r--	root/root	usr/lib32/pkgconfig/avahi-compat-libdns_sd.pc
-rw-r--r--	root/root	usr/lib32/pkgconfig/avahi-core.pc
-rw-r--r--	root/root	usr/lib32/pkgconfig/avahi-glib.pc
-rw-r--r--	root/root	usr/lib32/pkgconfig/avahi-gobject.pc
-rw-r--r--	root/root	usr/lib32/pkgconfig/avahi-libevent.pc
