# ports/core/mkinitcpio/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/initcpio/
drwxr-xr-x	root/root	etc/initcpio/hooks/
drwxr-xr-x	root/root	etc/initcpio/install/
drwxr-xr-x	root/root	etc/initcpio/post/
-rw-r--r--	root/root	etc/mkinitcpio.conf
drwxr-xr-x	root/root	etc/mkinitcpio.conf.d/
drwxr-xr-x	root/root	etc/mkinitcpio.d/
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/lsinitcpio
-rwxr-xr-x	root/root	usr/bin/mkinitcpio
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/initcpio/
-rwxr-xr-x	root/root	usr/lib/initcpio/functions
drwxrwxr-x	root/root	usr/lib/initcpio/hooks/
-rw-rw-r--	root/root	usr/lib/initcpio/hooks/consolefont
-rw-rw-r--	root/root	usr/lib/initcpio/hooks/encrypt
-rw-rw-r--	root/root	usr/lib/initcpio/hooks/keymap
-rw-rw-r--	root/root	usr/lib/initcpio/hooks/memdisk
-rw-rw-r--	root/root	usr/lib/initcpio/hooks/resume
-rw-rw-r--	root/root	usr/lib/initcpio/hooks/shutdown
-rw-rw-r--	root/root	usr/lib/initcpio/hooks/sleep
-rw-rw-r--	root/root	usr/lib/initcpio/hooks/udev
-rw-rw-r--	root/root	usr/lib/initcpio/hooks/usr
-rwxr-xr-x	root/root	usr/lib/initcpio/init
-rw-r--r--	root/root	usr/lib/initcpio/init_functions
drwxrwxr-x	root/root	usr/lib/initcpio/install/
-rw-rw-r--	root/root	usr/lib/initcpio/install/acpi_override
-rw-rw-r--	root/root	usr/lib/initcpio/install/autodetect
-rw-rw-r--	root/root	usr/lib/initcpio/install/base
-rw-rw-r--	root/root	usr/lib/initcpio/install/block
-rw-rw-r--	root/root	usr/lib/initcpio/install/consolefont
-rw-rw-r--	root/root	usr/lib/initcpio/install/encrypt
-rw-rw-r--	root/root	usr/lib/initcpio/install/filesystems
-rw-rw-r--	root/root	usr/lib/initcpio/install/fsck
-rw-rw-r--	root/root	usr/lib/initcpio/install/hostdata
-rw-rw-r--	root/root	usr/lib/initcpio/install/keyboard
-rw-rw-r--	root/root	usr/lib/initcpio/install/keymap
-rw-rw-r--	root/root	usr/lib/initcpio/install/kms
-rw-rw-r--	root/root	usr/lib/initcpio/install/lvm2
-rw-rw-r--	root/root	usr/lib/initcpio/install/mdadm_udev
-rw-rw-r--	root/root	usr/lib/initcpio/install/memdisk
-rw-rw-r--	root/root	usr/lib/initcpio/install/microcode
-rw-rw-r--	root/root	usr/lib/initcpio/install/modconf
-rw-rw-r--	root/root	usr/lib/initcpio/install/resume
-rw-rw-r--	root/root	usr/lib/initcpio/install/sd-encrypt
-rw-rw-r--	root/root	usr/lib/initcpio/install/sd-shutdown
-rw-rw-r--	root/root	usr/lib/initcpio/install/sd-vconsole
-rw-rw-r--	root/root	usr/lib/initcpio/install/shutdown
-rw-rw-r--	root/root	usr/lib/initcpio/install/sleep
-rw-rw-r--	root/root	usr/lib/initcpio/install/strip
-rw-rw-r--	root/root	usr/lib/initcpio/install/systemd
-rw-rw-r--	root/root	usr/lib/initcpio/install/udev
-rw-rw-r--	root/root	usr/lib/initcpio/install/usr
drwxr-xr-x	root/root	usr/lib/initcpio/post/
-rwxr-xr-x	root/root	usr/lib/initcpio/shutdown
drwxr-xr-x	root/root	usr/lib/initcpio/udev/
-rw-r--r--	root/root	usr/lib/initcpio/udev/01-memdisk.rules
drwxr-xr-x	root/root	usr/lib/kernel/
drwxr-xr-x	root/root	usr/lib/kernel/install.d/
-rwxr-xr-x	root/root	usr/lib/kernel/install.d/50-mkinitcpio.install
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
drwxr-xr-x	root/root	usr/lib/systemd/system/halt.target.wants/
lrwxrwxrwx	root/root	usr/lib/systemd/system/halt.target.wants/mkinitcpio-generate-shutdown-ramfs.service -> ../mkinitcpio-generate-shutdown-ramfs.service
drwxr-xr-x	root/root	usr/lib/systemd/system/kexec.target.wants/
lrwxrwxrwx	root/root	usr/lib/systemd/system/kexec.target.wants/mkinitcpio-generate-shutdown-ramfs.service -> ../mkinitcpio-generate-shutdown-ramfs.service
-rw-r--r--	root/root	usr/lib/systemd/system/mkinitcpio-generate-shutdown-ramfs.service
drwxr-xr-x	root/root	usr/lib/systemd/system/poweroff.target.wants/
lrwxrwxrwx	root/root	usr/lib/systemd/system/poweroff.target.wants/mkinitcpio-generate-shutdown-ramfs.service -> ../mkinitcpio-generate-shutdown-ramfs.service
drwxr-xr-x	root/root	usr/lib/systemd/system/reboot.target.wants/
lrwxrwxrwx	root/root	usr/lib/systemd/system/reboot.target.wants/mkinitcpio-generate-shutdown-ramfs.service -> ../mkinitcpio-generate-shutdown-ramfs.service
drwxr-xr-x	root/root	usr/lib/tmpfiles.d/
-rw-r--r--	root/root	usr/lib/tmpfiles.d/mkinitcpio.conf
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/bash-completion/
drwxr-xr-x	root/root	usr/share/bash-completion/completions/
lrwxrwxrwx	root/root	usr/share/bash-completion/completions/lsinitcpio -> mkinitcpio
-rw-r--r--	root/root	usr/share/bash-completion/completions/mkinitcpio
drwxr-xr-x	root/root	usr/share/libalpm/
drwxr-xr-x	root/root	usr/share/libalpm/hooks/
-rw-r--r--	root/root	usr/share/libalpm/hooks/60-mkinitcpio-remove.hook
-rw-r--r--	root/root	usr/share/libalpm/hooks/90-mkinitcpio-install.hook
drwxr-xr-x	root/root	usr/share/libalpm/scripts/
-rwxr-xr-x	root/root	usr/share/libalpm/scripts/mkinitcpio
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/lsinitcpio.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/mkinitcpio.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/mkinitcpio.8.gz
drwxr-xr-x	root/root	usr/share/mkinitcpio/
-rw-r--r--	root/root	usr/share/mkinitcpio/example.preset
-rw-r--r--	root/root	usr/share/mkinitcpio/hook.preset
drwxr-xr-x	root/root	usr/share/zsh/
drwxr-xr-x	root/root	usr/share/zsh/site-functions/
-rw-r--r--	root/root	usr/share/zsh/site-functions/_mkinitcpio
